# projects/mini-mail/container-fs/etc/s6/services/syslogd/run
#!/usr/bin/env bash

logger -p mail.info "s6-supervise : syslogd -> spawning process"
exec syslogd -n
