# projects/mini-mail/container-fs/etc/s6/services/dovecot/run
#!/usr/bin/env bash

logger -p mail.info "s6-supervise : dovecot -> spawning process"
dovecot -c /etc/dovecot/dovecot.conf -F
