# projects/frrouting/container-fs/etc/s6/services/snmpd/run
#!/usr/bin/env bash

/usr/sbin/snmpd -f -Le -LS0-6d -c /etc/snmp/snmpd.conf